# cron/crontab
*/1 * * * * curl web:80/app/track-show/process_clicks.php
*/1 * * * * curl web:80/app/track-show/process_postback.php
