# cron/crontab
*/1 * * * * curl -s web:80/track-show/process_clicks.php
*/1 * * * * curl -s web:80/track-show/process_postback.php
#* * * * * root echo "test cron..." >> /var/log/cron.log 2>&1